# workspace/crontab/www-data
* * * * * php /var/www/gutfan.com/web/artisan schedule:run >> /dev/null 2>&1
